# private_Library/LaunchAgents/dev.fzymgc.1password-agent-proxy.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>dev.fzymgc.1password-agent-proxy</string>
    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/socat</string>
        <string>UNIX-LISTEN:/Users/sean/.1password/agent.sock,fork,reuseaddr,unlink-early,mode=0600</string>
        <string>UNIX-CONNECT:/Users/sean/Library/Group Containers/2BUA8C4S2C.com.1password/t/agent.sock</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardErrorPath</key>
    <string>/tmp/1password-agent-proxy.err</string>
</dict>
</plist>
